# src/customfancontrol.service
[Unit]
Description=custom fan control
After=syslogd.target

[Service]
ExecStart=/usr/bin/customfancontrol

[Install]
WantedBy=multi-user.target

